# Freya/freyastrap/Library/dpkg/info/openssh-server.list
/.
/Library
/Library/LaunchDaemons
/Library/LaunchDaemons/com.openssh.sshd.plist
/etc
/etc/pam.d
/etc/pam.d/sshd
/etc/ssh
/etc/ssh/moduli
/etc/ssh/sshd_config
/usr
/usr/libexec
/usr/libexec/sshd-keygen-wrapper
/usr/sbin
/usr/sbin/sshd
/usr/share
/usr/share/doc
/usr/share/doc/openssh-server
/usr/share/doc/openssh-server/README.zst
/usr/share/doc/openssh-server/TODO.zst
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/moduli.5.zst
/usr/share/man/man5/sshd_config.5.zst
/usr/share/man/man8
/usr/share/man/man8/sshd.8.zst
