# zuesstrapNutzSigned/Library/dpkg/info/cydia.list
/.
/Applications
/Applications/Cydia.app
/Applications/Cydia.app/Cydia
/Applications/Cydia.app/Default-568h@2x.png
/Applications/Cydia.app/Default-Landscape.png
/Applications/Cydia.app/Default-Landscape@2x.png
/Applications/Cydia.app/Default-Portrait.png
/Applications/Cydia.app/Default-Portrait@2x.png
/Applications/Cydia.app/Default.png
/Applications/Cydia.app/Default@2x.png
/Applications/Cydia.app/Icon-60.png
/Applications/Cydia.app/Icon-60@2x.png
/Applications/Cydia.app/Icon-60@3x.png
/Applications/Cydia.app/Icon-72.png
/Applications/Cydia.app/Icon-72@2x.png
/Applications/Cydia.app/Icon-76@2x~ipad.png
/Applications/Cydia.app/Icon-76~ipad.png
/Applications/Cydia.app/Icon-83.5@2x.png
/Applications/Cydia.app/Icon-Small-40.png
/Applications/Cydia.app/Icon-Small-40@2x.png
/Applications/Cydia.app/Icon-Small-40@3x.png
/Applications/Cydia.app/Icon-Small-50.png
/Applications/Cydia.app/Icon-Small-50@2x.png
/Applications/Cydia.app/Icon-Small.png
/Applications/Cydia.app/Icon-Small@2x.png
/Applications/Cydia.app/Icon-Small@3x.png
/Applications/Cydia.app/Icon.png
/Applications/Cydia.app/Icon@2x.png
/Applications/Cydia.app/Info.plist
/Applications/Cydia.app/Sections.plist
/Applications/Cydia.app/changes.png
/Applications/Cydia.app/changes7.png
/Applications/Cydia.app/changes7@2x.png
/Applications/Cydia.app/changes7@3x.png
/Applications/Cydia.app/changes7s.png
/Applications/Cydia.app/changes7s@2x.png
/Applications/Cydia.app/changes7s@3x.png
/Applications/Cydia.app/changes@2x.png
/Applications/Cydia.app/chevron@2x.png
/Applications/Cydia.app/compose.png
/Applications/Cydia.app/configure.png
/Applications/Cydia.app/error.html
/Applications/Cydia.app/folder.png
/Applications/Cydia.app/folder@2x.png
/Applications/Cydia.app/home-Selected.png
/Applications/Cydia.app/home-Selected@2x.png
/Applications/Cydia.app/home.png
/Applications/Cydia.app/home7.png
/Applications/Cydia.app/home7@2x.png
/Applications/Cydia.app/home7@3x.png
/Applications/Cydia.app/home7s.png
/Applications/Cydia.app/home7s@2x.png
/Applications/Cydia.app/home7s@3x.png
/Applications/Cydia.app/home@2x.png
/Applications/Cydia.app/iOS7-Default-568h@2x.png
/Applications/Cydia.app/iOS7-Default-667h@2x.png
/Applications/Cydia.app/iOS7-Default-736h-Landscape@3x.png
/Applications/Cydia.app/iOS7-Default-736h-Portrait@3x.png
/Applications/Cydia.app/iOS7-Default-812h-Portrait@3x.png
/Applications/Cydia.app/iOS7-Default-Landscape.png
/Applications/Cydia.app/iOS7-Default-Landscape@2x.png
/Applications/Cydia.app/iOS7-Default-Portrait.png
/Applications/Cydia.app/iOS7-Default-Portrait@2x.png
/Applications/Cydia.app/iOS7-Default@2x.png
/Applications/Cydia.app/iconClassic.png
/Applications/Cydia.app/install.png
/Applications/Cydia.app/install7.png
/Applications/Cydia.app/install7@2x.png
/Applications/Cydia.app/install7@3x.png
/Applications/Cydia.app/install7s.png
/Applications/Cydia.app/install7s@2x.png
/Applications/Cydia.app/install7s@3x.png
/Applications/Cydia.app/install@2x.png
/Applications/Cydia.app/installed.png
/Applications/Cydia.app/installed@2x.png
/Applications/Cydia.app/installing.png
/Applications/Cydia.app/installing@2x.png
/Applications/Cydia.app/localize.js
/Applications/Cydia.app/manage.png
/Applications/Cydia.app/manage7.png
/Applications/Cydia.app/manage7@2x.png
/Applications/Cydia.app/manage7@3x.png
/Applications/Cydia.app/manage7s.png
/Applications/Cydia.app/manage7s@2x.png
/Applications/Cydia.app/manage7s@3x.png
/Applications/Cydia.app/manage@2x.png
/Applications/Cydia.app/reload.png
/Applications/Cydia.app/removing.png
/Applications/Cydia.app/removing@2x.png
/Applications/Cydia.app/search.png
/Applications/Cydia.app/search7.png
/Applications/Cydia.app/search7@2x.png
/Applications/Cydia.app/search7@3x.png
/Applications/Cydia.app/search7s.png
/Applications/Cydia.app/search7s@2x.png
/Applications/Cydia.app/search7s@3x.png
/Applications/Cydia.app/search@2x.png
/Applications/Cydia.app/unknown.png
/Applications/Cydia.app/Purposes
/Applications/Cydia.app/Purposes/commercial.png
/Applications/Cydia.app/Purposes/console.png
/Applications/Cydia.app/Purposes/daemon.png
/Applications/Cydia.app/Purposes/extension.png
/Applications/Cydia.app/Purposes/library.png
/Applications/Cydia.app/Purposes/uikit.png
/Applications/Cydia.app/Purposes/x.png
/Applications/Cydia.app/Sections
/Applications/Cydia.app/Sections/Addons.png
/Applications/Cydia.app/Sections/Administration.png
/Applications/Cydia.app/Sections/Archiving.png
/Applications/Cydia.app/Sections/Books.png
/Applications/Cydia.app/Sections/Carrier_Bundles.png
/Applications/Cydia.app/Sections/Data_Storage.png
/Applications/Cydia.app/Sections/Development.png
/Applications/Cydia.app/Sections/Dictionaries.png
/Applications/Cydia.app/Sections/Education.png
/Applications/Cydia.app/Sections/Entertainment.png
/Applications/Cydia.app/Sections/Fonts.png
/Applications/Cydia.app/Sections/Games.png
/Applications/Cydia.app/Sections/Health_and_Fitness.png
/Applications/Cydia.app/Sections/Java.png
/Applications/Cydia.app/Sections/Keyboards.png
/Applications/Cydia.app/Sections/Localization.png
/Applications/Cydia.app/Sections/Messaging.png
/Applications/Cydia.app/Sections/Multimedia.png
/Applications/Cydia.app/Sections/Navigation.png
/Applications/Cydia.app/Sections/Networking.png
/Applications/Cydia.app/Sections/Packaging.png
/Applications/Cydia.app/Sections/Productivity.png
/Applications/Cydia.app/Sections/Repositories.png
/Applications/Cydia.app/Sections/Ringtones.png
/Applications/Cydia.app/Sections/Scripting.png
/Applications/Cydia.app/Sections/Security.png
/Applications/Cydia.app/Sections/Site-Specific_Apps.png
/Applications/Cydia.app/Sections/Social.png
/Applications/Cydia.app/Sections/Soundboards.png
/Applications/Cydia.app/Sections/System.png
/Applications/Cydia.app/Sections/Terminal_Support.png
/Applications/Cydia.app/Sections/Text_Editors.png
/Applications/Cydia.app/Sections/Themes.png
/Applications/Cydia.app/Sections/Toys.png
/Applications/Cydia.app/Sections/Tweaks.png
/Applications/Cydia.app/Sections/Utilities.png
/Applications/Cydia.app/Sections/Wallpaper.png
/Applications/Cydia.app/Sections/Widgets.png
/Applications/Cydia.app/Sections/X_Window.png
/Applications/Cydia.app/Sources
/Applications/Cydia.app/_CodeSignature
/Applications/Cydia.app/_CodeSignature/CodeResources
/Applications/Cydia.app/menes
/Applications/Cydia.app/menes/menes.js
/Library
/Library/LaunchDaemons
/Library/LaunchDaemons/com.saurik.Cydia.Startup.plist
/etc
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/preferences.d
/etc/apt/preferences.d/cydia
/etc/apt/sources.list.d
/etc/apt/sources.list.d/cydia.list
/etc/apt/trusted.gpg.d
/etc/apt/trusted.gpg.d/bigboss.gpg
/etc/apt/trusted.gpg.d/hbang.gpg
/etc/apt/trusted.gpg.d/modmyi.gpg
/etc/apt/trusted.gpg.d/saurik.gpg
/etc/apt/trusted.gpg.d/sbingner.gpg
/etc/apt/trusted.gpg.d/zodttd.gpg
/usr
/usr/libexec
/usr/libexec/cydia
/usr/libexec/cydia/asuser
/usr/libexec/cydia/cfversion
/usr/libexec/cydia/cydo
/usr/libexec/cydia/du
/usr/libexec/cydia/finish.sh
/usr/libexec/cydia/firmware.sh
/usr/libexec/cydia/free.sh
/usr/libexec/cydia/move.sh
/usr/libexec/cydia/setnsfpn
/usr/libexec/cydia/startup
/var
/var/lib
/var/lib/cydia
/Applications/Cydia.app/bzip2
/Applications/Cydia.app/copy
/Applications/Cydia.app/file
/Applications/Cydia.app/gpgv
/Applications/Cydia.app/gzip
/Applications/Cydia.app/http
/Applications/Cydia.app/https
/Applications/Cydia.app/lzma
/Applications/Cydia.app/rred
/Applications/Cydia.app/store
/Applications/Cydia.app/Sections/Planet-iPhones Mods.png
/Applications/Cydia.app/Sources/apt.bigboss.us.com.png
