# extra tar junk/zuesstrapNutzSigned/private/etc/profile.d/cacerts.bootstrap.sh
## git ##
git config --global http.sslCAInfo /etc/ssl/certs/cacert.pem &> /dev/null
